# dashboard.service
[Unit]
Description=Minecraft Server Monitor Dashboard
After=network.target docker.service

[Service]
Type=simple
WorkingDirectory=%h/dashboard
EnvironmentFile=-%h/dashboard/.env
ExecStart=/usr/bin/python3 %h/dashboard/dashboard.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
